# 99-mount-external-drive.rules
SUBSYSTEM=="block", ENV(ID_FS_UUID)=="0F1B-6ED2", ACTION=="add", RUN+="/usr/local/bin/mount_external_drive.sh"
